# notifier.service
[Unit]
Description=Gmail & Calendar Telegram Notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fahmi
ExecStart=/home/fahmi/bin/notifier
Environment=NOTIFIER_CONFIG=/home/fahmi/.notifier/config.yaml
Restart=no
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
